# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:telegram_bot]
command=python bot.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/bot_err.log
stdout_logfile=/var/log/supervisor/bot_out.log
environment=PYTHONUNBUFFERED=1

[program:web_server]
command=python web_server.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/webserver_err.log
stdout_logfile=/var/log/supervisor/webserver_out.log
environment=PYTHONUNBUFFERED=1
